# data/supervisor.conf
[program:autotracking]
user = pinachet
directory = /
command = uwsgi --emperor /etc/uwsgi/vassals 
autostart = true
autorestart = true
stderr_logfile = /var/www/AutoTracking/data/log/uwsgi-err.log
stdout_logfile = /var/www/AutoTracking/data/log/uwsgi-out.log
stopsignal = QUIT
